# systemd/vc.service
[Unit]
Description=VSR Technologies P25 Trunking (VC)
After=syslog.target zerotier-one.service network.target

[Service]
User=root
Type=forking
ExecStart=/opt/dvm/bin/dvmhost -c /opt/dvm/bin/VC.yml
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
